# templates/pssid-daemon.service.j2
[Unit]
Description=Pssid-Daemon

[Service]
Restart=always
WorkingDirectory=/usr/bin/pssid/pssid-daemon
ExecStart=/usr/bin/python3 pssid-daemon.py

[Install]
WantedBy=multi-user.target
